# configs/crontab
1 3,15 * * * python /home/ubuntu/VirtualEnvironment/news_classification/daily_job.py
